# scripts/com.balboa.morning-briefing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.balboa.morning-briefing</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/vimalkumarasamy/Workspace/lemonade-agent/.venv/bin/python</string>
        <string>/Users/vimalkumarasamy/Workspace/lemonade-agent/scripts/morning_briefing.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/vimalkumarasamy/Workspace/lemonade-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/vimalkumarasamy</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/balboa-briefing.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/balboa-briefing-error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
